# rootfs/etc/init.d/mihowrt
#!/bin/sh /etc/rc.common

START=21
STOP=79
USE_PROCD=1

SERVICE_NAME="mihowrt"
CLASH_DIR="/opt/clash"
CLASH_BIN="${CLASH_DIR}/bin/clash"
CLASH_CONFIG="${CLASH_DIR}/config.yaml"
ORCHESTRATOR="/usr/bin/mihowrt"
SKIP_START_FILE="/tmp/luci-app-mihowrt.skip-start"
SERVICE_PID_FILE="/var/run/mihowrt/mihomo.pid"
SERVICE_START_TIMEOUT="${SERVICE_START_TIMEOUT:-5}"
EXTRA_COMMANDS="apply"
EXTRA_HELP="        apply   Validate and apply on-disk MihoWRT config/policy"

msg() {
	logger -p daemon.info -t "$SERVICE_NAME" "$*"
}

service_running() {
	[ -x "$ORCHESTRATOR" ] || return 1
	"$ORCHESTRATOR" service-running >/dev/null 2>&1
}

validate_service_inputs() {
	[ -x "$CLASH_BIN" ] || {
		msg "ERROR: Mihomo binary not found at $CLASH_BIN"
		return 1
	}

	[ -f "$CLASH_CONFIG" ] || {
		msg "ERROR: Mihomo config not found at $CLASH_CONFIG"
		return 1
	}

	if ! "$CLASH_BIN" -d "$CLASH_DIR" -t; then
		msg "ERROR: Mihomo configuration test failed"
		return 1
	fi

	if ! "$ORCHESTRATOR" validate; then
		msg "ERROR: Policy validation failed"
		return 1
		fi
}

wait_for_service_running() {
	local waited=0
	local timeout="$SERVICE_START_TIMEOUT"

	while [ "$waited" -lt "$timeout" ]; do
		if service_running; then
			return 0
		fi

		sleep 1
		waited=$((waited + 1))
	done

	return 1
}

start_service() {
	[ -f "$SKIP_START_FILE" ] && {
		msg "Skipping MihoWRT service auto-start during installer transaction"
		return 0
	}

	msg "Starting MihoWRT service..."
	"$ORCHESTRATOR" recover >/dev/null 2>&1 || true

	if ! validate_service_inputs; then
		return 1
	fi

	if ! "$ORCHESTRATOR" cleanup >/dev/null 2>&1; then
		msg "ERROR: Failed to clean stale runtime state"
		return 1
	fi

	procd_open_instance
	procd_set_param command "$ORCHESTRATOR" run-service
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param respawn 3600 5 5
	procd_close_instance

	if ! wait_for_service_running; then
		stop >/dev/null 2>&1 || true
		msg "ERROR: MihoWRT service did not stay running after start"
		return 1
	fi

	msg "MihoWRT service registered with procd"
}

stop_service() {
	msg "Stopping MihoWRT service..."
}

reload_service() {
	msg "Reloading MihoWRT policy..."

	if ! service_running; then
		msg "MihoWRT service is not running; skipping policy reload"
		return 0
	fi

	if ! "$ORCHESTRATOR" reload-policy >/dev/null 2>&1; then
		msg "ERROR: Policy reload failed"
		return 1
	fi

	msg "MihoWRT policy reloaded"
}

apply() {
	msg "Applying MihoWRT on-disk config..."

	[ -f "$SKIP_START_FILE" ] && {
		msg "ERROR: Cannot apply while installer skip-start marker is active"
		return 1
	}

	if ! validate_service_inputs; then
		return 1
	fi

	if ! service_running; then
		msg "MihoWRT service is not running; validated on-disk config only"
		return 0
	fi

	msg "MihoWRT service is running; restarting to apply on-disk changes"

	if ! stop; then
		msg "ERROR: Failed to stop MihoWRT service for apply"
		return 1
	fi

	if ! start; then
		stop >/dev/null 2>&1 || true
		msg "ERROR: Failed to restart MihoWRT service with updated on-disk config"
		return 1
	fi

	msg "MihoWRT on-disk changes applied"
}

service_triggers() {
	procd_add_reload_trigger "mihowrt"
}
